# px3se_linux/device/rockchip/px3-se/cif-camera-udev/rules.d/rk-cif-camera.rules
ACTION=="change", ENV{CVBS_NAME}=="*RESOLUTION*", RUN+="/etc/rk-camera-uevent.sh $env{CVBS_NAME}"
